# roles/clickhouse_backup_trigger/templates/clickhouse-backup-trigger.timer.j2
[Unit]
Description=clickhouse-backup-trigger timer

[Timer]
Unit=clickhouse-backup-trigger.service
OnCalendar={{ clickhouse_backup_schedule }}

[Install]
WantedBy=timers.target
